# slack-master/trunk/src/slack-master.cron
# $Header$
# vim:filetype=crontab
#
# Update the slack-master source
*/10 * * * *    root    /usr/lib/slack-master/slack-update-source
