# slack-master/trunk/src/slack-master.cron
# $Header$
# vim:filetype=crontab
#
# Update the slack-master source
*/10 * * * *    root    test -x /usr/sbin/slack-update-source && /usr/sbin/slack-update-source
